# com.fuellog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fuellog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jmneal/fuel-log/run-server.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jmneal/fuel-log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/jmneal/fuel-log/server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jmneal/fuel-log/server.err.log</string>
</dict>
</plist>
